# etc/systemd/system/nginx.service.d/override.conf
[Service]
# With read-only mounted root file system, /var/log is tmpfs so we need to create the directory before starting nginx
LogsDirectory=nginx
